# app/deploy/gunicorn.service
[Unit]
Description=Gunicorn instance for FastAPI Course Management app
After=network.target

[Service]
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/app
Environment="PATH=/home/ec2-user/venv/bin"
Environment="PYTHONPATH=/home/ec2-user/app"

# AWS credentials should be managed through IAM roles or environment variables
Environment="AWS_ACCESS_KEY_ID="
Environment="AWS_SECRET_ACCESS_KEY="
Environment="AWS_REGION="
Environment="AWS_SESSION_TOKEN="

# Application configuration
Environment="DYNAMODB_TABLE=dalhousie-courses"
Environment="S3_BUCKET=dalhousie-course-files"

ExecStart=/home/ec2-user/venv/bin/gunicorn \
    -w 4 \
    -k uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --access-logfile /var/log/gunicorn/access.log \
    --error-logfile /var/log/gunicorn/error.log \
    --capture-output \
    --log-level info \
    main:app

[Install]
WantedBy=multi-user.target 